# ventoy/devboost-firstboot.service
[Unit]
Description=dev-boost first-boot bootstrap (devboost install full)
After=network-online.target
Wants=network-online.target
ConditionPathExists=/opt/dev-boost/devboost

[Service]
Type=oneshot
RemainAfterExit=yes
# Runs the full provisioning once on the first networked boot, then disables itself so
# subsequent boots are normal. The frozen binary is the entrypoint (no bash). Secrets are
# the age-encrypted bundle injected onto the USB.
Environment=DEVBOOST_SECRETS=/opt/dev-boost/secrets.age
ExecStart=/bin/sh -c '/opt/dev-boost/devboost install full >> /var/log/devboost-firstboot.log 2>&1; systemctl disable devboost-firstboot.service'

[Install]
WantedBy=multi-user.target
