# ventoy/devboost-firstboot.service
[Unit]
Description=dev-boost first-boot bootstrap (devboost install full)
After=network-online.target
Wants=network-online.target
ConditionPathExists=/opt/dev-boost/devboost

[Service]
Type=oneshot
RemainAfterExit=yes
# All three vars are required by the secrets module (CONTRACT: secrets.py reads
# DEVBOOST_BOOTSTRAP_DIR and must NOT fall back to ".").
Environment=DEVBOOST_BOOTSTRAP_DIR=/opt/dev-boost
Environment=DEVBOOST_SECRETS=/opt/dev-boost/secrets.age
Environment=DEVBOOST_SECRETS_KEY=/opt/dev-boost/age-key.txt
# Run the provisioning binary directly; redirect output to a persistent log.
ExecStart=/opt/dev-boost/devboost install full
StandardOutput=append:/var/log/devboost-firstboot.log
StandardError=append:/var/log/devboost-firstboot.log
# Self-disable ONLY on success (exit 0). On failure the service stays enabled so the
# operator can retry: `systemctl reset-failed devboost-firstboot && systemctl start devboost-firstboot`
ExecStartPost=/bin/systemctl disable devboost-firstboot.service

[Install]
WantedBy=multi-user.target
