# dotfiles/dot_config/systemd/user/browser-mcp.service
[Unit]
Description=Playwright MCP — per-session browsers for remote Claude Code (Tailscale)
Documentation=https://github.com/microsoft/playwright-mcp
# Needs the tailnet up (the launcher binds to the Tailscale IP) and, for --isolated
# Chromium, a graphical session to display the browsers.
After=tailscaled.service network-online.target graphical-session.target
Wants=network-online.target

[Service]
Type=simple
# The launcher auto-detects the Tailscale IP and shells out to `npx`, which mise
# provides via shims — put them on PATH since a user unit starts with a bare one.
Environment=PATH=%h/.local/share/mise/shims:%h/.local/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=%h/.local/bin/browser-mcp
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
